# build/service/coldfusion9/run
#!/bin/bash

set -e

exec /opt/coldfusion9/bin/coldfusion start
